# deploy/systemd/pubgd-api.service
[Unit]
Description=PUBG dashboard API (FastAPI/uvicorn)
Documentation=file:///home/pubg/pubg_dashboard/HANDOFF.md
After=network-online.target
Wants=network-online.target
# See pubgd-worker.service: a [Unit] directive, silently ignored under [Service].
StartLimitIntervalSec=0

[Service]
Type=exec
WorkingDirectory=/home/pubg/pubg_dashboard/backend
Environment=PATH=/home/pubg/.local/bin:/usr/local/bin:/usr/bin:/bin
# Bound to the LAN by explicit decision — there is no authentication, so
# anything that can reach this host can read the dashboard and can also POST to
# /api/players and /api/ingest, which mutate state and spend PUBG rate-limit
# budget. Postgres and MinIO are deliberately NOT exposed (see
# docker/docker-compose.yml); only this port is.
ExecStart=/home/pubg/.local/bin/uv --directory /home/pubg/pubg_dashboard/backend run \
    uvicorn pubg_dashboard.api.app:app --host 0.0.0.0 --port 8000 --log-level info

Restart=always
RestartSec=10s
KillSignal=SIGTERM
TimeoutStopSec=20s

# Read-only endpoints, but /replay reads whole bundles into memory to serve
# them; a few concurrent requests for the largest match is the worst case.
MemoryHigh=400M
MemoryMax=600M

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
